# docker/config/supervisord-php7-fpm.conf
[program:php7]
command=/usr/sbin/php-fpm7.2
