# docker/config/supervisord-php5.6-fpm.conf
[program:php5]
command=/usr/sbin/php-fpm5.6
